# etc/cron.hourly/update-store
#!/bin/sh
cd /var/www/item.tf
./updatestore.py
