# etc/cron.hourly/update-cache
#!/bin/sh
cd /var/www/item.tf
./updatecache.py
