# etc/cron.hourly/update-store
#!/bin/sh
cd /srv/www/item.tf
./updatestore.py
