# etc/rc.d/itemtf
#!/bin/sh

# PROVIDE: itemtf
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name=itemtf

load_rc_config $name
: ${itemtf_user:=item.tf}
: ${itemtf_chdir:=/usr/local/www/item.tf}
: ${itemtf_prepend:="/usr/sbin/daemon -f"}
: ${itemtf_env:="PYTHONOPTIMIZE=2"}

command=/usr/local/libexec/item.tf/bin/uvicorn
command_args=main:app
command_interpreter=/usr/local/libexec/item.tf/bin/python3.10

run_rc_command "$1"
