# etc/cron.daily/backup
#!/bin/sh
cd /var/www/item.tf
mkdir -p backups
cp /var/lib/redis/dump.rdb backups/$(date +%Y%m%d%H%M).rdb
